# docker/nginx.sv.conf
[program:webserver]
command=/bin/bash -c "/usr/sbin/nginx -g 'daemon off;'"
redirect_stderr=true
