# docker/tomcat.sv.conf
[program:webapp]
command=/bin/bash -c "catalina.sh run"
redirect_stderr=true
